# deploy/aprilcamd.service
[Unit]
Description=AprilCam tag-detection daemon
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=eric
SupplementaryGroups=video

# systemd creates and manages these directories before ExecStart runs.
# Paths resolved to:
#   ConfigurationDirectory → /etc/aprilcam   (or ~/.config/aprilcam when --user)
#   StateDirectory         → /var/lib/aprilcam
#   RuntimeDirectory       → /run/aprilcam
#   LogsDirectory          → /var/log/aprilcam
ConfigurationDirectory=aprilcam
StateDirectory=aprilcam
RuntimeDirectory=aprilcam
LogsDirectory=aprilcam

# APRILCAM_DATA_DIR points to the user's persistent data (cameras, playfields).
# APRILCAM_SOCKET_DIR is the runtime dir managed by systemd (ephemeral).
Environment=APRILCAM_DATA_DIR=/home/eric/aprilcam-data
Environment=APRILCAM_SOCKET_DIR=/run/aprilcam

# Run the daemon in the foreground so systemd captures output via journald.
# --tcp   : bind gRPC on all interfaces (port 5280) for remote MCP clients.
# --unix  : also bind a Unix socket under APRILCAM_SOCKET_DIR for on-box CLI.
ExecStart=/home/eric/.local/share/pipx/venvs/aprilcam/bin/python \
    -m aprilcam.daemon \
    --tcp --tcp-port 5280 \
    --unix

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
